# assets/noctalia-auth.service.in
[Unit]
Description=Noctalia Auth - Unified Authentication Agent
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
ExecStartPre=@LIBEXECDIR@/noctalia-auth-bootstrap
ExecStartPre=/bin/sh -c 'mkdir -p "$${XDG_DATA_HOME:-$$HOME/.local/share}/dbus-1/services" && install -m 0644 @DATADIR@/noctalia-auth/org.gnome.keyring.SystemPrompter.service "$${XDG_DATA_HOME:-$$HOME/.local/share}/dbus-1/services/org.gnome.keyring.SystemPrompter.service"'
ExecStart=@LIBEXECDIR@/noctalia-auth --daemon
Slice=session.slice
TimeoutStopSec=5sec
Restart=on-failure

# Security hardening
PrivateTmp=yes
ProtectSystem=strict
NoNewPrivileges=true
CapabilityBoundingSet=
SystemCallFilter=@system-service
SystemCallArchitectures=native
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=true

[Install]
WantedBy=graphical-session.target
